# cups-airprint/rootfs/etc/s6-overlay/s6-rc.d/dbus-daemon/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting D-Bus daemon from S6"

# Garante que o diretório existe
mkdir -p /var/run/dbus

# Remove PID file antigo se existir
rm -f /var/run/dbus/pid

exec dbus-daemon --system --nofork
