# cups-airprint/rootfs/etc/s6-overlay/s6-rc.d/initialization/up
#!/usr/bin/env bashio

# Prepara diretórios necessários para o addon
mkdir -p /var/run/avahi-daemon
mkdir -p /run/cups

